# deploy/systemd/miniden-bot.service
[Unit]
Description=MiniDeN Telegram bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/miniden
EnvironmentFile=/opt/miniden/.env
ExecStart=/opt/miniden/venv/bin/python bot.py
Restart=on-failure
RestartSec=5s
User=miniden
Group=miniden

[Install]
WantedBy=multi-user.target
