# deploy/systemd/miniden-bot.service
[Unit]
Description=MiniDeN Telegram bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/miniden
EnvironmentFile=/opt/miniden/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/miniden/venv/bin/python bot.py
User=miniden
Group=miniden
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
