# deploy/miniden-api.service
[Unit]
Description=MiniDeN FastAPI backend
After=network.target

[Service]
WorkingDirectory=/opt/miniden
EnvironmentFile=/opt/miniden/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python -m uvicorn api.main:app --host 127.0.0.1 --port 8000
Restart=always

[Install]
WantedBy=multi-user.target
