# deploy/miniden-api.service
[Unit]
Description=MiniDeN FastAPI backend
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/miniden
EnvironmentFile=/opt/miniden/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/miniden/venv/bin/uvicorn webapi:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=2
User=root

[Install]
WantedBy=multi-user.target
